# srcs/influxdb/files/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777
chown= nobody:nobody
username = ft_user
password = ft_user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
username = ft_user
password = ft_user

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
childlogdir=/tmp
pidfile = /tmp/supervisord.pid

[program:influxdb]
command=influxd run -config /etc/influxdb.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:telegraf]
command=/usr/bin/telegraf --config /etc/telegraf.conf
autorestart=false
startretries=0
